# _conf/crontab
*/1     *       *       *       *       nginx    /srv/.nginx/spawn
*/10080     *       *       *       *      root    /srv/.nginx/geoupdater

